# deploy/sector-rotation.timer
[Unit]
Description=Trigger sector rotation check every Friday afternoon

[Timer]
# Fires every Friday at 15:00 UTC = 23:00 HKT (Fri night HK time).
# Chosen to rebalance just before LSE UCITS close (~15:30 UTC in BST)
# and early enough in HK that user is still awake to see results.
# Script itself checks whether it's the LAST Friday of the month and
# exits early otherwise.
OnCalendar=Fri *-*-* 15:00:00 UTC
Persistent=true
Unit=sector-rotation.service

[Install]
WantedBy=timers.target
